# docker/supervisor/queue-worker.conf
; Supervisor configuration for the fleet-arrival queue workers.
;
; Two worker pools process arrivals in parallel so independent planet destinations
; never wait for each other. The per-destination lock (stored in the shared
; `database` cache) keeps same-planet processing correct no matter how many workers run.
;
;   light pool  -> drains ONLY the light lane (transports, deployments, returns, ...).
;                  Battles are never queued here, so light traffic always has workers.
;   heavy pool  -> drains the heavy lane (attacks, moon destruction) first, then
;                  backfills the light lane when no battles are pending.
;
; This guarantees a simple transport is never stuck behind a storm of large battles.
; The {{QUEUE_WORKERS_LIGHT}} / {{QUEUE_WORKERS_HEAVY}} placeholders are substituted at
; container start from the matching environment variables (see docker/entrypoint.sh).

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:fleet-queue-worker-light]
process_name=%(program_name)s_%(process_num)02d
; --sleep=0.5 keeps delayed fleet arrivals snappy (default sleep is 3s).
; --timeout=600 matches ProcessFleetArrival::$timeout so large batches are not killed mid-job.
command=php /var/www/artisan queue:work --queue=fleet-arrivals --sleep=0.5 --timeout=600 --no-interaction
numprocs={{QUEUE_WORKERS_LIGHT}}
autostart=true
autorestart=true
; Give a worker up to 610s (> the 600s job timeout) to finish its current job on graceful stop.
stopwaitsecs=610
stopsignal=TERM
user=www-data
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:fleet-queue-worker-heavy]
process_name=%(program_name)s_%(process_num)02d
; Prioritise the heavy lane, then help drain the light lane when no battles are pending.
command=php /var/www/artisan queue:work --queue=fleet-arrivals-heavy,fleet-arrivals --sleep=0.5 --timeout=600 --no-interaction
numprocs={{QUEUE_WORKERS_HEAVY}}
autostart=true
autorestart=true
stopwaitsecs=610
stopsignal=TERM
user=www-data
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
